# systemd/content.service.in
[Unit]
Description = Service to download and make rss items readable

[Service]
Type = oneshot
ExecStart = BIN_NAME -path DATA_PATH
