# systemd/dispatch.service.in
[Unit]
Description = Service to dispatch readable articles

[Service]
Type = oneshot
ExecStart = BIN_NAME -path DATA_PATH
