# systemd/feeds.service.in
[Unit]
Description = Service to check new items in rss feeds

[Service]
Type = oneshot
ExecStart = BIN_NAME --path DATA_PATH
