# units/content.timer
[Unit]
Description=Run service every couple of hours

[Timer]
RandomizedDelaySec=20
OnCalendar=*-*-* *:*/5:00

[Install]
WantedBy=timers.target
